# app/supervisord.conf
[supervisord]
nodaemon=true

[program:festivalapi]
command=gunicorn wsgi --config /festival/gunicorn.conf --chdir=/festival
